# scripts/launchd/local.claude-code-local-proxy.plist.in
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>${AUTOSTART_LABEL}</string>
	<key>ProgramArguments</key>
	<array>
		<string>${UV_BIN}</string>
		<string>run</string>
		<string>claude-code-local-proxy</string>
		<string>--log-file</string>
		<string>${PROXY_LOG_FILE}</string>
	</array>
	<key>WorkingDirectory</key>
	<string>${WORKING_DIRECTORY}</string>
	<key>EnvironmentVariables</key>
	<dict>
		<key>PATH</key>
		<string>${LAUNCHD_PATH}</string>
	</dict>
	<key>RunAtLoad</key>
	<true/>
	<key>KeepAlive</key>
	<dict>
		<key>SuccessfulExit</key>
		<false/>
	</dict>
	<key>StandardOutPath</key>
	<string>${STDOUT_LOG_FILE}</string>
	<key>StandardErrorPath</key>
	<string>${STDERR_LOG_FILE}</string>
</dict>
</plist>
